# crontab
#CO2
@reboot /home/pi/CO2.sh
*/10 *  *   *   *     /home/pi/CO2.sh
